# dot_config/systemd/user/wallpaper-rotate.service
[Unit]
Description=Rotate wallpaper from ~/Pictures/Wallpapers

[Service]
Type=oneshot
ExecStart=%h/.config/hypr/custom/scripts/wallpaper-rotate.sh
